# host/.old/rubicon/Library/LaunchAgents/com.harmsway.populate_notes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harmsway.populate_notes</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/dharms/bin/populate_notes.sh</string>
  </array>

  <key>Nice</key>
  <integer>1</integer>

  <key>StartCalendarInterval</key>
  <dict>
	<key>Minute</key>
	<integer>0</integer>
   </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/populate_notes.err</string>

  <key>StandardOutPath</key>
  <string>/tmp/populate_notes.out</string>
</dict>
</plist>
